# systemd/deal-aggregator.timer
# Generated by scripts/generate_systemd_units.py — do not edit by hand.
# Translated from launchd plist: com.greenwich-barrow.deal-aggregator

[Unit]
Description=G&B deal-aggregator timer

[Timer]
OnCalendar=Mon..Fri 06:00:00
Persistent=true
Unit=deal-aggregator.service

[Install]
WantedBy=timers.target
