# systemd/deal-aggregator-afternoon.timer
# Generated by scripts/generate_systemd_units.py — do not edit by hand.
# Translated from launchd plist: com.greenwich-barrow.deal-aggregator-afternoon

[Unit]
Description=G&B deal-aggregator-afternoon timer

[Timer]
OnCalendar=Mon..Fri 14:00:00
Persistent=true
Unit=deal-aggregator-afternoon.service

[Install]
WantedBy=timers.target
