# systemd/email-intelligence.timer
# Generated by scripts/generate_systemd_units.py — do not edit by hand.
# Translated from launchd plist: com.greenwich-barrow.email-intelligence

[Unit]
Description=G&B email-intelligence timer

[Timer]
OnCalendar=Mon..Fri 07:00:00
Persistent=true
Unit=email-intelligence.service

[Install]
WantedBy=timers.target
