# systemd/health-monitor.timer
# Generated by scripts/generate_systemd_units.py — do not edit by hand.
# Translated from launchd plist: com.greenwich-barrow.health-monitor

[Unit]
Description=G&B health-monitor timer

[Timer]
OnCalendar=Fri 00:30:00
Persistent=true
Unit=health-monitor.service

[Install]
WantedBy=timers.target
