# systemd/jj-snapshot-refresh.timer
# Generated by scripts/generate_systemd_units.py — do not edit by hand.
# Translated from launchd plist: com.greenwich-barrow.jj-snapshot-refresh

[Unit]
Description=G&B jj-snapshot-refresh timer

[Timer]
OnCalendar=Mon..Fri 09:00:00
OnCalendar=Mon..Fri 14:30:00
OnCalendar=Mon..Fri 18:00:00
Persistent=true
Unit=jj-snapshot-refresh.service

[Install]
WantedBy=timers.target
