# systemd/launchd-debugger.timer
# Generated by scripts/generate_systemd_units.py — do not edit by hand.
# Translated from launchd plist: com.greenwich-barrow.launchd-debugger

[Unit]
Description=G&B launchd-debugger timer

[Timer]
OnCalendar=*-*-* 05:00:00
Persistent=true
Unit=launchd-debugger.service

[Install]
WantedBy=timers.target
