# systemd/niche-intelligence.timer
# Generated by scripts/generate_systemd_units.py — do not edit by hand.
# Translated from launchd plist: com.greenwich-barrow.niche-intelligence

[Unit]
Description=G&B niche-intelligence timer

[Timer]
OnCalendar=Tue 22:30:00
Persistent=true
Unit=niche-intelligence.service

[Install]
WantedBy=timers.target
