# systemd/nightly-tracker-audit.timer
# Generated by scripts/generate_systemd_units.py — do not edit by hand.
# Translated from launchd plist: com.greenwich-barrow.nightly-tracker-audit

[Unit]
Description=G&B nightly-tracker-audit timer

[Timer]
OnCalendar=*-*-* 23:30:00
Persistent=true
Unit=nightly-tracker-audit.service

[Install]
WantedBy=timers.target
