# systemd/post-call-analyzer-poll.service
# Hand-authored 2026-05-08 (Phase 4). NOT generated by
# scripts/generate_systemd_units.py — there is no source plist on the iMac
# because this is the new server-side MCP-based detector that replaces
# post_call_analyzer_poll.py during shadow mode. The iMac sidecar plist
# (com.greenwich-barrow.post-call-analyzer.plist) keeps running unchanged
# until shadow-mode validation completes (Phase 4.5).
#
# Pairs with: post-call-analyzer-poll.timer

[Unit]
Description=G&B post-call-analyzer MCP poll (Phase 4 server-side detector)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/projects/Sapling
EnvironmentFile=%h/projects/Sapling/scripts/.env.launchd
Environment=PATH=%h/.npm-global/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 %h/projects/Sapling/scripts/post_call_analyzer_mcp_poll.py

# Cap any one tick — protects the timer from a hung claude -p subprocess
# that somehow escaped the script's 180s internal timeout.
TimeoutStartSec=300

[Install]
WantedBy=default.target
