# systemd/post-call-analyzer-poll.timer
# Server-side timer for post-call-analyzer detector. Rewritten 2026-05-13.
#
# Cadence: 2 fires/day at 1pm + 6pm America/New_York (handles DST automatically).
# Persistent=true so a missed fire after VPS reboot triggers once on wake.
#
# Per session-decisions 2026-05-13: Kay decided 5min was too aggressive given
# 1-3 calls/day. Midday + EOD is the right cadence — first fire catches morning
# calls in time to surface in afternoon, EOD fire catches afternoon calls + any
# late-arriving transcript updates.

[Unit]
Description=G&B post-call-analyzer poll timer (1pm + 6pm ET)

[Timer]
OnCalendar=*-*-* 13:00:00 America/New_York
OnCalendar=*-*-* 18:00:00 America/New_York
Persistent=true
Unit=post-call-analyzer-poll.service

[Install]
WantedBy=timers.target
