# systemd/post-call-analyzer-poll.service
# Server-side detector for post-call-analyzer. Rewritten 2026-05-13 from
# MCP-based Python poller to REST-based Bash poller via granola-api wrapper.
# No MCP, no OAuth, no reconnect. Credentials resolve from 1Password.
#
# Pairs with: post-call-analyzer-poll.timer

[Unit]
Description=G&B post-call-analyzer detector (REST via granola-api)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/projects/Sapling
EnvironmentFile=%h/projects/Sapling/scripts/.env.launchd
Environment=PATH=%h/.npm-global/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/bin/bash %h/projects/Sapling/scripts/post_call_analyzer_poll.sh

# Cap any one tick — protects timer from a hung subprocess
TimeoutStartSec=600

[Install]
WantedBy=default.target
