# systemd/relationship-manager.timer
# Generated by scripts/generate_systemd_units.py — do not edit by hand.
# Translated from launchd plist: com.greenwich-barrow.relationship-manager

[Unit]
Description=G&B relationship-manager timer

[Timer]
OnCalendar=Mon..Fri 06:50:00
Persistent=true
Unit=relationship-manager.service

[Install]
WantedBy=timers.target
