# systemd/weekly-snapshot.timer
# Generated by scripts/generate_systemd_units.py — do not edit by hand.
# Translated from launchd plist: com.greenwich-barrow.weekly-snapshot

[Unit]
Description=G&B weekly-snapshot timer

[Timer]
OnCalendar=Fri 22:00:00
Persistent=true
Unit=weekly-snapshot.service

[Install]
WantedBy=timers.target
