# systemd/weekly-archive-export.timer
# Generated by scripts/generate_systemd_units.py — do not edit by hand.
# Translated from launchd plist: com.greenwich-barrow.weekly-archive-export

[Unit]
Description=G&B weekly-archive-export timer

[Timer]
OnCalendar=Sat 09:00:00
Persistent=true
Unit=weekly-archive-export.service

[Install]
WantedBy=timers.target
